# modules/openwrt-mihomo/files/etc/uci-defaults/90-mihomo
uci batch << EOF
set mihomo.subscription.url='$CLASH_CONFIG_URL'
commit
EOF