# modules/opkg-mirror/files/etc/uci-defaults/99-opkg
sed -e 's,https://downloads.immortalwrt.org,https://mirror.nju.edu.cn/immortalwrt,g' \
    -e 's,https://mirrors.vsean.net/openwrt,https://mirror.nju.edu.cn/immortalwrt,g' \
    -i.bak /etc/opkg/distfeeds.conf