# modules/network/files/etc/uci-defaults/89-network
uci batch << EOF
set network.wan.proto='pppoe'
set network.wan.username='$PPPOE_USERNAME'
set network.wan.password='$PPPOE_PASSWORD'

commit
EOF

uci batch << EOF

set network.lan.ipaddr='$LAN_IP'

commit
EOF