# modules/ssh-permission/files/etc/uci-defaults/99-ssh
if [ -f /etc/dropbear/authorized_key ]; then
    chmod 644 /etc/dropbear/authorized_keys
fi