# modules/prefer-ipv6-settings/files/etc/uci-defaults/99-ipv6
uci batch << EOF
delete network.globals.ula_prefix
del dhcp.lan
set dhcp.lan=dhcp
set dhcp.lan.interface='lan'
set dhcp.lan.start='100'
set dhcp.lan.limit='150'
set dhcp.lan.leasetime='12h'
set dhcp.lan.dhcpv4='server'
set dhcp.lan.ra='server'
set dhcp.lan.ra_useleasetime='1'
set dhcp.lan.preferred_lifetime='8h'
set dhcp.lan.ra_flags='none'
set dhcp.lan.dns_service='0'

commit
EOF