# modules/root-password/files/etc/uci-defaults/91-system
(echo "$ROOT_PASSWORD"; sleep 1; echo "$ROOT_PASSWORD") | passwd > /dev/null