# modules/reject-netflix-ipv6/files/etc/uci-defaults/99-reject-netflix-ipv6
echo "conf-dir=/etc/dnsmasq.d" >> /etc/dnsmasq.conf