# modules/daed-as-default/files/etc/uci-defaults/99-daed
uci del dhcp.@dnsmasq[0].dns_redirect
uci set daed.config.enabled='1'
uci commit